# etc/systemd/system/hblock.service.d/override.conf
[Service]
Environment=HBLOCK_WRAP=20
Environment=HBLOCK_REGEX=true